# packer/files/etc/rc.local
#!/bin/sh -e

iptables -t nat -A POSTROUTING ! -d 10.0.0.0/8 -o eth0 -j MASQUERADE
